# package/rooter/0routerspecfic/h721/files/etc/uci-defaults/999-force
#!/bin/sh

uci set dhcp.lan.force='1'
uci commit dhcp
/etc/init.d/dnsmasq restart
#rm -f /etc/uci-defaults/999-force
#/usr/lib/custom/reboot.sh